# nippelboard.service
[Unit]
Description=Nippelboard Funk - Radio Monitoring & Soundboard
After=network.target sound.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Nippelboard_Funk
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/usr/bin/python3 /home/pi/Nippelboard_Funk/main.py
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nippelboard

# Security (optional)
# PrivateTmp=yes
# NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
